# docker/icommands/supervisord.conf
[supervisord]
nodaemon=true

[program:postgres]
command=/etc/init.d/postgresql start

[program:irods]
command=/etc/init.d/irods start