# docker/icat/supervisord.conf
[supervisord]
nodaemon=true

[program:irods]
command=/etc/init.d/irods start